# linux/calibre-upgrade.timer
[Unit]
Description=Timer to upgrade calibre

[Timer]
OnCalendar=Fri 06:00
Persistent=True

[Install]
WantedBy=multi-user.target
